# launchagents/com.recall.auto-index.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.recall.auto-index</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{RECALL_MCP_DIR}}/scripts/run-auto-index.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>{{RECALL_MCP_DIR}}/logs/auto-index-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>{{RECALL_MCP_DIR}}/logs/auto-index-launchd-error.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
